# graph-generator/sv/autoregenerate-sleep-graphs/run
#!/bin/sh
exec /usr/local/bin/run_when_changed.sh \
	watch="$DATA_DIR" \
	files="sleep.csv" \
	script="/usr/local/bin/on_sleep.csv_changed.sh"

